# Library/LaunchAgents/com.yqrashawn.mbsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yqrashawn.mbsync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/watch</string>
    <string>--no-title</string>
    <string>--beep</string>
    <string>-n 240</string>
    <string>proxychains4 -f /etc/proxychains.conf mbsync gmail</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>