# Library/LaunchAgents/com.yqrashawn.dark-mode-notify.plist
<?xml version="1.0" encoding="utf-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.yqrashawn.dark-mode-notify</string>
    <key>KeepAlive</key>
    <true />
    <key>StandardErrorPath</key>
    <string>/tmp/dark-mode-notify-stderr.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/dark-mode-notify-stdout.log</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/dark-mode-notify</string>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>exec
      $HOME/local/bin/dark_light_mode_change.clj</string>
    </array>
  </dict>
</plist>
